# discovery-server/ros2-discovery-server.service
[Unit]
Description=ROS2 Jazzy Fast DDS Discovery Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Environment="ROS_DISTRO=jazzy"

# Use the actual binary - the fastdds wrapper is broken
ExecStart=/bin/bash -c 'source /opt/ros/jazzy/setup.bash && fast-discovery-server -i 0 -p 11811'

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ros2-discovery-server

[Install]
WantedBy=multi-user.target
